# aports/nodm/nodm.initd
#!/sbin/openrc-run

depends() {
	needs localmount dbus
}

start() {
	ebegin "Starting nodm"

	NODM_USER=${NODM_USER} NODM_XSESSION=${NODM_XSESSION} start-stop-daemon --start -b -m -p /var/run/nodm.pid --exec /usr/bin/nodm
	eend $?
}

stop() {
	ebegin "Stopping nodm"
	start-stop-daemon --stop -p /var/run/nodm.pid --exec /usr/bin/nodm
	eend $?
}
